# dist/supervisor/entree.conf
[program:entree]
command = /home/cs411/.virtualenvs/entree/bin/gunicorn_start          ; Command to start app
user = cs411                                                          ; User to run as
stdout_logfile = /home/cs411/entree/dist/logs/gunicorn_supervisor.log ; Where to write log messages
redirect_stderr = true                                                ; Save stderr in the same log
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8                       ; Set UTF-8 as default encoding